# example/supervisord.conf
; =============================================================================
;  zzci/supervisord — full-feature reference configuration
; =============================================================================
;
; Every section and every per-program key documented below is recognised by
; the daemon. Drop-in compatible with Python supervisord for the parts that
; are still implemented.
;
; Sections that are PARSED BUT IGNORED in this fork (logged once at startup
; with a warning):
;   - [inet_http_server]                    -> unix-socket-only listener
;   - [unix_http_server].username/password  -> filesystem permissions are auth
;
; Removed entirely vs. upstream:
;   - Web GUI, JSON REST handlers, BasicAuth + SHA1 path
;
; =============================================================================


; -----------------------------------------------------------------------------
;  [unix_http_server] — XML-RPC listener (the one that's used)
; -----------------------------------------------------------------------------
[unix_http_server]
file = /var/run/supervisor.sock     ; absolute path to the unix socket
; chmod                             ; not implemented (use umask + mount)
; chown                             ; not implemented
; username / password               ; ignored with a warning

; The TCP listener is parsed but ignored. Keep this section out unless you
; have a reason to keep upstream-compat config files round-trippable.
; [inet_http_server]
; port = 127.0.0.1:9001
; username = admin
; password = secret


; -----------------------------------------------------------------------------
;  [supervisord] — daemon-wide settings
; -----------------------------------------------------------------------------
[supervisord]
logfile          = /var/log/supervisord.log
logfile_maxbytes = 50MB              ; rotate after this size; suffixes K/M/G/KB/MB/GB
logfile_backups  = 10                ; number of rotated files to keep
loglevel         = info              ; trace | debug | info | warn | error | critical
pidfile          = /var/run/supervisord.pid
identifier       = supervisor        ; used in XML-RPC supervisor.getIdentification
minfds           = 1024              ; raise RLIMIT_NOFILE soft to at least this; aborts on failure
minprocs         = 200               ; raise RLIMIT_NPROC soft to at least this


; -----------------------------------------------------------------------------
;  [supervisorctl] — defaults read by `supervisord ctl`
; -----------------------------------------------------------------------------
[supervisorctl]
serverurl = unix:///var/run/supervisor.sock
; username / password are accepted but ignored on the unix socket; auth is
; the socket's filesystem permissions.


; -----------------------------------------------------------------------------
;  [include] — pull additional .conf fragments from disk (glob expanded)
; -----------------------------------------------------------------------------
; [include]
; files = /etc/supervisord.d/*.conf  conf.d/*.conf


; =============================================================================
;  Program reference — every recognised key
; =============================================================================
;
; This block names every per-program option the daemon honours. Real
; deployments will only set a subset; the smaller examples further down show
; common patterns.
[program:reference]
; --- core --------------------------------------------------------------------
command          = /usr/bin/myserver --port=8080      ; required; argv parsed shell-style
process_name     = %(program_name)s                   ; default; can use %(process_num)s in templates
numprocs         = 1                                  ; spawn N replicas of this program
numprocs_start   = 0                                  ; first process_num value (default 0)
priority         = 999                                ; lower = earlier in the start order
autostart        = true                               ; start on supervisord boot
autorestart      = unexpected                         ; true | false | unexpected (only on non-listed exit codes)
startsecs        = 1                                  ; needed alive this long to enter Running
startretries     = 3                                  ; consecutive start failures before Fatal
restartpause     = 0                                  ; seconds to sleep between restart attempts
exitcodes        = 0,2                                ; comma-separated "expected" exit codes
stopsignal       = TERM                               ; signal(s) to gracefully stop, space-separated for retry chain
stopwaitsecs     = 10                                 ; seconds to wait between stopsignal sends
killwaitsecs     = 2                                  ; seconds before SIGKILL after the final stopsignal
stopasgroup      = false                              ; send signal to whole process group
killasgroup      = false                              ; SIGKILL to whole process group; defaults to stopasgroup
user             = www-data                           ; setuid+setgid before exec; supports user[:group]
directory        = /srv/myserver                      ; chdir before exec
environment      = HTTP_ADDR=":8080",DEBUG="0"        ; literal env vars
envFiles         = /etc/myserver.env                  ; comma-separated dotenv files
redirect_stderr  = false                              ; merge child stderr into stdout

; --- stdout / stderr capture -------------------------------------------------
stdout_logfile          = /var/log/myserver.out.log    ; AUTO | NONE | path
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups  = 10
stdout_capture_maxbytes = 0                            ; <PROCESS_COMMUNICATION> capture buffer
stdout_events_enabled   = false                        ; emit PROCESS_LOG_STDOUT events
stderr_logfile          = /var/log/myserver.err.log
stderr_logfile_maxbytes = 50MB
stderr_logfile_backups  = 10
stderr_capture_maxbytes = 0
stderr_events_enabled   = false

; --- syslog (optional alternative to file logging) --------------------------
; syslog_facility       = local6
; syslog_tag            = myserver
; syslog_stdout_priority= info
; syslog_stderr_priority= warning

; --- cron-style scheduling (alternative to autostart) -----------------------
; cron = 0 */5 * * * *                                 ; runs the command on this cron expr (6-field, with seconds)

; --- restart on file/binary change ------------------------------------------
restart_when_binary_changed         = false
restart_cmd_when_binary_changed     =                 ; if set, run instead of stop+start
restart_signal_when_binary_changed  =                 ; if set, send signal instead

restart_directory_monitor           =                 ; absolute dir path to watch (file change monitor)
restart_file_pattern                = *
restart_cmd_when_file_changed       =
restart_signal_when_file_changed    =

; --- depends_on (real readiness gating) -------------------------------------
depends_on          = db, cache                       ; comma-separated program names
depends_on_timeout  = 60                              ; seconds to wait per dep before strategy kicks in
depends_on_strategy = abort                           ; abort (default) | ignore
depends_on_ready    = healthy                         ; running | healthy

; --- healthchecks (set any one to enable; ALL must pass per attempt) --------
healthcheck_http             = http://127.0.0.1:8080/healthz
healthcheck_tcp              = 127.0.0.1:8080
healthcheck_file             = /run/myserver.ready
healthcheck_command          = curl -fsS http://localhost:8080/healthz
healthcheck_timeout          = 5                      ; per-attempt timeout (seconds)
healthcheck_interval         = 2                      ; gap between attempts (seconds)
healthcheck_retries          = 3                      ; consecutive successes to flip Healthy
healthcheck_startup_timeout  = 300                    ; ceiling on the wait-for-Running phase


; =============================================================================
;  Realistic snippets — copy/paste starting points
; =============================================================================

; ---- HTTP service with a /healthz endpoint ---------------------------------
[program:web]
command         = /usr/bin/myserver --addr=:8080
autostart       = true
autorestart     = true
startsecs       = 2
stopsignal      = TERM
stopwaitsecs    = 15
stdout_logfile  = /var/log/web.out.log
stderr_logfile  = /var/log/web.err.log
healthcheck_http     = http://127.0.0.1:8080/healthz
healthcheck_interval = 5
healthcheck_retries  = 3

; ---- Database needing a TCP-port healthcheck -------------------------------
[program:db]
command         = /usr/bin/postgres -D /var/lib/postgres
user            = postgres
autostart       = true
autorestart     = true
startsecs       = 5
stopsignal      = INT TERM           ; pg accepts SIGINT for fast shutdown, fall back to TERM
stopwaitsecs    = 30
healthcheck_tcp     = 127.0.0.1:5432
healthcheck_retries = 5

; ---- Background worker that depends on db being healthy --------------------
[program:worker]
command          = /usr/bin/myworker
autostart        = true
autorestart      = true
startsecs        = 3
depends_on       = db
depends_on_ready = healthy
depends_on_timeout = 120

; ---- Cron-driven script (one-shot per fire, not autostart) -----------------
[program:nightly-cleanup]
command   = /usr/local/bin/cleanup.sh
autostart = false                    ; the cron expression starts it
cron      = 0 0 3 * * *              ; sec min hour dom mon dow — 03:00 daily

; ---- Group: stop/start as a unit -------------------------------------------
[group:stack]
programs = db,worker,web

; =============================================================================
;  Event listeners — hook into supervisord's event protocol
; =============================================================================
;
; Same surface as [program:*] plus an `events` key. Communication is over
; the listener's own stdin/stdout (PROCESS_COMMUNICATION).
;
; [eventlistener:state-monitor]
; command       = /usr/local/bin/notify-on-fatal.py
; events        = PROCESS_STATE_FATAL,PROCESS_STATE_BACKOFF
; buffer_size   = 100
; autostart     = true
; stdout_logfile= /var/log/state-monitor.out.log
